# crontab
MAILTO=francis@flourish.org
@hourly ./tool/scraper.js > /dev/null
